# scripts/webapp.service
[Unit]
Description=CSYE 6225 Web Application
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=csye6225
Group=csye6225
EnvironmentFile=/etc/environment
WorkingDirectory=/opt/csye6225
ExecStart=/usr/bin/java -jar /opt/csye6225/webapp.jar
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webapp

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
